# files/etc/systemd/system/skyhub-startup.service
[Unit]
Description=Sky Hub Startup
After=multi-user.target

[Service]
ExecStart=/bin/sh /bin/skyhub-startup

[Install]
WantedBy=multi-user.target

