# files/etc/cron.daily/skyhub-jobs.sh
#!/bin/bash

skyhub purge
skyhub upgrade
