# files/skyhub-pipe.service
[Unit]
Description=Sky Hub command named pipe

[Service]
Type=simple
ExecStart=/bin/bash /bin/skyhub-pipe

[Install]
WantedBy=multi-user.target

